# recipes-core/firsttime-boot/files/firsttime-boot.target
[Unit]
Description=First time boot
Conflicts=rescue.service rescue.target
Before=local-fs-pre.target
ConditionFirstBoot=yes
AllowIsolate=no
